# templates/default/systemd/alertmanager.service.erb
[Unit]
Description=Prometheus Alertmanager
After=network.target

[Service]
ExecStart=<%= node['prometheus']['alertmanager']['binary'] %> \
    -log.level=debug \
    -storage.path='/var/lib/prometheus/data' \
    -config.file='/etc/prometheus/alertmanager.conf' \
    -web.external-url=<%= node['prometheus']['alertmanager.url'] %>
User=<%= node['prometheus']['user'] %>
Group=<%= node['prometheus']['group'] %>
Restart=always

[Install]
WantedBy=multi-user.target
